# deploy/systemd/agent-sozluk-runtime.service
[Unit]
Description=Agent Sozluk PostgreSQL queue runtime orchestrator
Documentation=file:/opt/agent-sozluk/runtime/current/docs/PRODUCTION_RUNBOOK.md
Wants=network-online.target
After=network-online.target
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=agent-runtime
Group=agent-runtime
WorkingDirectory=/opt/agent-sozluk/runtime/current
EnvironmentFile=/etc/agent-sozluk/runtime.env
Environment=NODE_ENV=production
Environment=NO_COLOR=1
Environment=PATH=/usr/bin:/usr/local/bin:/bin

# Fail closed if the operator installed a different Node major or omitted a
# required host artifact. pnpm resolves the versioned project-local tsx binary.
ExecStartPre=/usr/bin/node -e "process.exit(Number(process.versions.node.split('.')[0]) === 22 ? 0 : 1)"
ExecStartPre=/usr/bin/pnpm --version
ExecStartPre=/usr/bin/test -x /usr/local/bin/codex
ExecStartPre=/usr/bin/test -x /usr/bin/bwrap
ExecStartPre=/usr/bin/bwrap --version
ExecStartPre=/usr/bin/test -r /var/lib/agent-sozluk-runtime/credentials.json
ExecStart=/usr/bin/pnpm exec tsx scripts/agent-runtime-worker.ts

Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
KillMode=mixed
TimeoutStopSec=45s
SendSIGKILL=yes
FinalKillSignal=SIGKILL

UMask=0077
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes
DevicePolicy=closed
ProtectSystem=strict
ProtectHome=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectProc=invisible
ReadOnlyPaths=/opt/agent-sozluk/runtime/current /var/lib/agent-sozluk-runtime/credentials.json /etc/agent-sozluk/runtime.env
ReadWritePaths=/opt/agent-sozluk/runtime/codex-home /opt/agent-sozluk/runtime/work
InaccessiblePaths=-/opt/agent-sozluk/app -/run/docker.sock -/var/run/docker.sock
CapabilityBoundingSet=
AmbientCapabilities=
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=user mnt pid ipc uts
SystemCallArchitectures=native
# Bubblewrap needs user/mount/PID namespaces and mount syscalls to construct
# the per-child view. It receives no capability and cannot mount on the host.
SystemCallFilter=~@clock @cpu-emulation @debug @module @obsolete @raw-io @reboot @swap

MemoryHigh=1536M
MemoryMax=2048M
TasksMax=128
LimitNOFILE=4096
CPUQuota=200%
OOMPolicy=stop
Nice=5

StandardOutput=journal
StandardError=journal
SyslogIdentifier=agent-sozluk-runtime
LogRateLimitIntervalSec=30s
LogRateLimitBurst=200

[Install]
WantedBy=multi-user.target
